# files/etc/crontabs/root
* * * * * /sbin/test_gateway > /dev/null
* * * * * /sbin/splash_sync > /dev/null
