# files/etc/crontabs/root
* * * * * /sbin/test_gateway > /dev/null
