# files/common/etc/hotplug.d/iface/22-wan_fsm
#!/bin/sh
[ "$INTERFACE" = "wan" ] && {
	logger -t hotplug.d "Wan interface changed state, calling fsm watch inetable"
        /sbin/fsm watch inetable
}
