# files/common/etc/crontabs/root
* * * * * /sbin/p2ptbl gossip /tmp/p2ptbl/test 100 br-mesh
