# files/common/etc/crontabs/root
* * * * */10 /sbin/p2ptbl gossip /tmp/p2ptbl/ccfg_wifi 500 br-mesh
*/1 * * * *    /sbin/p2ptbl gossip /tmp/p2ptbl/gwip      500 br-mesh
*/1 * * * *    /sbin/p2ptbl gossip /tmp/p2ptbl/splash    500 br-mesh
*/5 * * * *    /sbin/p2ptbl gossip /tmp/p2ptbl/update    500 br-mesh
*/1 * * * * /sbin/fsm watch inetable
*/5 * * * * /sbin/fsm watch update
*/1 * * * * /sbin/splash_sync
*/5 * * * * /sbin/n2n_watchdog
