# files/common/etc/crontabs/root
* * * * * /sbin/p2ptbl gossip /tmp/p2ptbl/update 500 br-mesh
* * * * * /sbin/fsm watch update
