# files/common/etc/init.d/hbbp
#!/bin/sh /etc/rc.common
START=45

start() {
	cd /hbbp
	listener br-mesh
}

stop() {
       killall listener
}
